# scripts/pinsheet.service
[Unit]
Description=PinSheet Modern
After=network.target

[Service]
Type=simple
User=pinsheet
WorkingDirectory=/opt/pinsheet
Environment=SECRET_KEY=REPLACE_ME
Environment=FLASK_DEBUG=0
ExecStart=/opt/pinsheet/.venv/bin/python source/main.py --host 127.0.0.1 --port 8080
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
