# scripts/pinsheet.service
[Unit]
Description=PinSheet Server
After=network.target

[Service]
Type=simple
User=hickory
WorkingDirectory=/opt/pinsheet-server
Environment=SECRET_KEY=REPLACE_ME
Environment=FLASK_ENV=production
ExecStart=/opt/pinsheet-server/.venv/bin/python source/main.py --host 0.0.0.0 --port 8080
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
